# monitor_hotplug/etc/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl restart monitor-hotplug"
